# ansible/files/etc/cron.d/update-tor-blocklist
0 5 * * * root /usr/local/bin/update-tor-blocklist.sh
